# Docker/render-mongo/supervisord.conf
[supervisord]
nodaemon=true 

[program:mongod]
command=/usr/bin/mongod --dbpath /data/mongo

[program:render]
command=/bin/bash -c "${JETTY_HOME}/bin/jetty.sh run"

